# service/bambu-discord-notifier.service
[Unit]
Description=Bambu Lab Discord Notifier
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/BambuDiscordNotifier
ExecStart=/home/pi/BambuDiscordNotifier/venv/bin/python run.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
